# luci-app-simplenetwork/root/etc/uci-defaults/40_luci-app-simplenetwork

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@simplenetwork[-1]
	add ucitrack simplenetwork
	set ucitrack.@simplenetwork[-1].init=simplenetwork
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
