# luci-app-simplenetwork/root/etc/uci-defaults/40_luci-app-simplenetwork
touch /etc/config/simplenetwork
touch /etc/config/network
touch /etc/config/wireless
touch /etc/config/firewall
touch /etc/config/dhcp

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@simplenetwork[-1]
	add ucitrack simplenetwork
	set ucitrack.@simplenetwork[-1].init=simplenetwork
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

# system/system
sed -i 's,"path": "system/system","path": "system/system-simple",' /usr/share/luci/menu.d/luci-mod-system.json

# system/flash
sed -i 's,"path": "system/flash","path": "system/flash-simple",' /usr/share/luci/menu.d/luci-mod-system.json

# hide ssh
sed -i 's,"luci-mod-system-ssh","luci-mod-system-ssh-hidden",' /usr/share/luci/menu.d/luci-mod-system.json
sed -i 's,"luci-mod-system-init","luci-mod-system-init-hidden",' /usr/share/luci/menu.d/luci-mod-system.json
sed -i 's,"luci-mod-system-cron","luci-mod-system-cron-hidden",' /usr/share/luci/menu.d/luci-mod-system.json
sed -i 's,"luci-mod-system-mounts","luci-mod-system-mounts-hidden",' /usr/share/luci/menu.d/luci-mod-system.json

# hide network
sed -i 's,"admin/network","admin/network-hidden",' /usr/share/luci/menu.d/luci-base.json

exit 0
