# luci-app-simplenetwork/root/etc/uci-defaults/40_luci-app-simplenetwork
touch /etc/config/simplenetwork
touch /etc/config/network
touch /etc/config/wireless
touch /etc/config/firewall
touch /etc/config/dhcp

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@simplenetwork[-1]
	add ucitrack simplenetwork
	set ucitrack.@simplenetwork[-1].init=simplenetwork
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache

exit 0
