# luci-app-wizard/files/root/etc/init.d/wizard
#!/bin/sh /etc/rc.common

START=19

boot() {
	:
}

add_wizard() {
	:
}

start() {
	config_load wizard
	config_foreach add_wizard wizard
}

restart() {
	XRELOAD=1 start
}
