# net/mwan3plus/files/etc/uci-defaults/mwan3-reload-firewall
#!/bin/sh

uci -q batch <<-EOT
	delete firewall.mwan3_ipv6_masq
	delete firewall.mwan3
	commit firewall
EOT

exit 0
