# net/mwan3plus/files/etc/uci-defaults/mwan3-ipv6-masq-help

uci -q batch <<-EOT
	delete firewall.mwan3_ipv6_masq
	set firewall.mwan3_ipv6_masq=include
	set firewall.mwan3_ipv6_masq.type=script
	set firewall.mwan3_ipv6_masq.path=/lib/mwan3/mwan3_ipv6_masq.sh
	set firewall.mwan3_ipv6_masq.family=ipv6
	set firewall.mwan3_ipv6_masq.reload=0
	commit firewall
EOT

exit 0
