# natcap/files/natcapd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=95

start() {
	if test $(ls -l /data/crashdump/ 2>/dev/null | wc -l) -gt 0 ; then
		touch /tmp/pstore
	fi
	SERVICE_DAEMONIZE=1 service_start /usr/sbin/natcapd start
	echo start ok
}

stop() {
	/usr/sbin/natcapd stop
	echo stop ok
}

restart() {
	echo restart
	start
}

reload() {
	echo reload
	start
}
