# natcap/files/natcapd_boot.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=10

boot() {
	test -f /etc/config/natcap_server.sh && sh /etc/config/natcap_server.sh
	/usr/sbin/natcapd boot
}
