# base-config-setting/files/disk_ready_done.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org
# Copyright (C) 2019 X-WRT <dev@x-wrt.com>

START=96

boot() {
	sync
	mount | grep overlayfs | grep -q 'workdir=/overlay/work' && {
		test -f /rom/etc/sda.ready || touch /rom/etc/sda.ready
		mount -o remount,ro /rom
		sync
	}
}
