# base-config-setting/files/ethtool.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=55

boot() {
	for IFACE in `cat /proc/net/dev | grep eth | awk -F: '{print $1}'`; do
		ethtool -K $IFACE gro off >/dev/null 2>&1
	done
}
