# base-config-setting/files/sta_disable.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org
# Copyright (C) 2019 X-WRT <dev@x-wrt.com>

START=99

start() {
	SERVICE_DAEMONIZE=1 service_start /usr/sbin/sta_disable
}

stop() {
	:
}

restart() {
	:
}

reload() {
	:
}
