# telnetd-test/files/telnetd-test.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=98

start() {
	telnetd=$(fw_printenv telnetd 2>/dev/null | cut -d= -f2)
	[ "$telnetd" = "off" ] && {
		/etc/init.d/telnetd-test disable
		return 0
	}
	
	/usr/sbin/telnetd -l /bin/sh &
}

stop() {
	killall telnetd
}
