# luci-app-fakemesh/root/etc/init.d/fakemesh
#!/bin/sh /etc/rc.common

START=99

boot() {
	XBOOT=1 start
}

start() {
	/usr/sbin/fakemesh
}

restart() {
	XRELOAD=1 start
}
